# backups/coolify_recovery_20250624_021755/backup/crontab
# Sauvegardes automatiques MCP
# Tous les jours à 2h00
0 2 * * * /usr/local/bin/backup.sh

# Nettoyage des anciennes sauvegardes tous les dimanche à 3h00  
0 3 * * 0 /usr/local/bin/cleanup.sh 